# Scripts/autorun/com.dwellium.frontend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dwellium.frontend</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/ilyaklipinitser/Library/Application Support/Dwellium/start-frontend.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ilyaklipinitser/Library/Application Support/Dwellium/frontend-build</string>

    <key>RunAtLoad</key>
    <true/>

    <!-- Auto-restart on crash (non-zero exit OR crash signal) -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/ilyaklipinitser/Library/Logs/dwellium-frontend.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/ilyaklipinitser/Library/Logs/dwellium-frontend.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>production</string>
        <key>HOME</key>
        <string>/Users/ilyaklipinitser</string>
        <!-- Default port; can be overridden by editing this plist + reloading -->
        <key>DWELLIUM_FRONTEND_PORT</key>
        <string>5173</string>
    </dict>

    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
